# usr/lib/systemd/system/android-mount.service.d/99-cpuset.conf
[Service]
Environment="BIND_MOUNT_PATH=/var/lib/lxc/android/rootfs"
ExecStart=
ExecStart=/usr/sbin/mount-android.sh
