# usr/lib/systemd/user/miatoll-quirks.service
[Unit]
Description=Miatoll quirks

[Service]
ExecStart=/usr/bin/miatoll-quirks.sh

[Install]
WantedBy=default.target
