# usr/lib/systemd/system/android-cpuset.service
[Unit]
Description=Setup android cpuset

[Service]
Type=oneshot
ExecStart=/usr/sbin/cpuset_init.sh

[Install]
WantedBy=sysinit.target
