# usr/lib/systemd/user/org.gnome.SettingsDaemon.MediaKeys.service.d/99-cpuset.conf
[Service]
ExecStartPost=+/usr/bin/cpuset_add.sh %n background
