# usr/lib/systemd/system/systemd-resolved.service.d/99-cpuset.conf
[Service]
ExecStartPost=+/usr/sbin/cpuset_add.sh %n system-background
